# com.neetware.discord-agent-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.neetware.discord-agent-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jackkim/Projects/discord-agent-bridge/run-launchagent.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/jackkim/Projects/discord-agent-bridge</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/jackkim/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/jackkim/Library/Logs/discord-agent-bridge/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jackkim/Library/Logs/discord-agent-bridge/stderr.log</string>
</dict>
</plist>
